# images/precise/base/provision/etc/sv/sshd/run
#!/bin/bash
set -e
rm -f /etc/ssh/*key*
ssh-keygen -A
mkdir -p /var/run/sshd
exec /usr/sbin/sshd -De -p $SSHD_PORT 2>&1
